# infra/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisor/supervisord.pid

; ---- Redis process ----
[program:redis]
command=redis-server
    --appendonly yes
    --requirepass %(ENV_REDIS_PASSWORD)s
    --bind 127.0.0.1
    --port 6379
    --dir /data
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ---- oracle-sim Node.js process ----
; Starts after Redis (priority=20) so Redis is ready to accept connections.
[program:oracle-sim]
command=node /app/packages/oracle-sim/dist/server.js
autostart=true
autorestart=true
priority=20
; REDIS_HOST is localhost because both processes share the same container
environment=
    NODE_ENV="production",
    REDIS_HOST="127.0.0.1",
    REDIS_PORT="6379",
    REDIS_PASSWORD="%(ENV_REDIS_PASSWORD)s",
    ORACLE_PORT="8787"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
